# .config/systemd/user/restic-cleanup.timer
[Unit]
Description=Purge backups once a week

[Timer]
OnCalendar=Sun *-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target
