# .config/systemd/user/restic.service
[Unit]
Description=Run restic backup to backblaze b2

[Service]
Type=oneshot
ExecStart=/usr/bin/restic backup %h --exclude-file %h/.resticignore
